# src/backend/services/sync/crontab
# START CRON JOB
* * * * * /app/src/backend/services/sync/publisher.sh > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB